# .dotfiles/profile.d/brew.pyenv
#!/usr/bin/env sh

if brew_keg_exists pyenv; then
  eval "$(pyenv init -)"
fi
