# .dotfiles/profile.d/brew.sqlite
#!/usr/bin/env sh

if brew_keg_exists sqlite; then
  prepend_to_path /usr/local/opt/sqlite/bin
fi
