# .dotfiles/profile.d/brew-cask.atom
#!/usr/bin/env sh

if command_exists atom
then
  alias aatom='atom -a'
  export GRAPHICAL_EDITOR=atom
fi
