# .dotfiles/profile.d/shell
#!/usr/bin/env sh

# Sets tabs to 4 chars long
tabs -4
